# supervisord-kivitask.conf
[program:task_server.pl]
command=/usr/local/bin/supervisord-kivitask.sh
numprocs=1
autostart=true
autorestart=true
startsecs=5
priority=10
exitcodes=0
