# supervisord-cups-browsed.conf
[program:cups-browsed]
command=/usr/local/bin/supervisord-cups-browsed.sh
numprocs=1
autostart=true
autorestart=true
startsecs=8
priority=10
exitcodes=0

