# supervisord-exim4.conf
[program:exim4]
command=/usr/local/bin/supervisord-exim4.sh
numprocs=1
autostart=true
autorestart=true
startsecs=8
priority=10
exitcodes=0

