# supervisord-cups.conf
[program:cups]
command=/usr/local/bin/supervisord-cups.sh
numprocs=1
autostart=true
autorestart=true
startsecs=8
priority=10
exitcodes=0

